# files/check_openmanage.cron
#!/bin/bash

# Dynamically generated by Puppet
# Local modifications will be over-written

# Run check_openmanage
/usr/local/bin/check_openmanage -all &> /var/log/check_openmanage.out
